# images/php/supervisord/template/consumer.ini
[program:%event%]
directory=/var/www/html/apiapps
process_name=%event%
command=php artisan queue:work %event% --sleep=3 --tries=3 --daemon
autostart=true
autorestart=true
startretries=999999
startsecs=5
redirect_stderr=false
stdout_logfile=/var/www/html/apiapps/var/logs/%event%.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/var/www/html/apiapps/var/logs/%event%.error.log
stderr_capture_maxbytes=1MB
